# src/apps/daemon/service/dde-cooperation-daemon.service
[Unit]
Description=Deepin Cooperation Daemon
After=network.target

[Service]
ExecStart=/usr/bin/dde-cooperation-daemon
Restart=always
RestartSec=10
RestartPreventExitStatus=SIGTERM

[Install]
WantedBy=default.target
